# kafka-role/templates/zookeeper.service.j2
[Unit]
Description=zookeeper
After=syslog.target network.target

[Service]
Type=simple

User={{ kafka_user }}
Group={{ kafka_group }}
Environment=KAFKA_LOG4J_OPTS="-Dlog4j.configuration=file:{{ log4j_dir }}"
Environment=KAFKA_HEAP_OPTS="-Xmx1G -Xms1G"

ExecStart={{ kafka_bin_dir }}/zookeeper-server-start.sh {{ zookeeper_conf_file }}
ExecStop={{ kafka_bin_dir }}/zookeeper-server-stop.sh

[Install]
WantedBy=multi-user.target
